# jobs/mail_delivery/move_to_pipeline.timer
[Unit]
Description=Run pipeline delivery job every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
